# Vault/vault_a.conf
[program:vault]
command=/usr/local/bin/vault server -config=vault-config-a.hcl
user=root
stderr_logfile=/var/log/vault.err.log
stdout_logfile=/var/log/vault.out.log
